# DEBIAN/postrm
#!/bin/bash

systemctl disable counter

systemctl daemon-reload